# usr/lib/systemd/system/systemd-repart.service.d/30_grub-live.conf
## Copyright (C) 2025 - 2025 ENCRYPTED SUPPORT LLC <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

[Unit]

## Debian forked upstream module 90overlay-root
ConditionKernelCommandLine=!rootovl

## dracut upstream module 90overlayfs - grub-live / ISO
ConditionKernelCommandLine=!rd.live.overlay.overlayfs=1

## initramfs-tools - live-boot
ConditionKernelCommandLine=!boot=live
ConditionKernelCommandLine=!nopersistence
